# dist/files/ovn-kubernetes-master.service
[Unit]
Description=ovn-kubernetes master service
Requires=openvswitch.service
After=openvswitch.service
Requires=ovn-northd.service
After=ovn-northd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/ovn-kubernetes
ExecStart=/usr/bin/ovn-kubernetes-master.sh

[Install]
WantedBy=multi-user.target

