# dist/files/ovn-kubernetes-node.service
[Unit]
Description=Provision ovn-kubernetes node service
Requires=openvswitch.service
After=openvswitch.service
Requires=ovn-controller.service
After=ovn-controller.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/ovn-kubernetes
ExecStart=/usr/bin/ovn-kubernetes-node.sh

[Install]
WantedBy=multi-user.target

